# deploy/systemd/ledger@.service
# Ledger Contábil — Systemd Template
[Unit]
Description=Ledger Contábil (%i)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/vivify/releases/%i/ledger
EnvironmentFile=/var/www/vivify/shared/.env
Environment=COLOR=%i
ExecStart=/var/www/vivify/deploy/bin/ledger-start %i

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
NoNewPrivileges=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_UNIX
ReadWritePaths=/var/www/vivify/shared /var/www/vivify/releases/%i/ledger
LimitNOFILE=16384
MemoryMax=256M
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
